# systemd/.config/systemd/user/cleanup-maildir.timer
[Unit]
Description = cleanup maildir folders daily

[Timer]
OnStartupSec = 1h
OnUnitActiveSec = 24h

[Install]
WantedBy = timers.target
