# systemd/.config/systemd/user/dir_size_monitor.timer
[Unit]
Description = monitor changes on size of dirs daily

[Timer]
OnCalendar = 20:45

[Install]
WantedBy = timers.target
