# systemd/.config/systemd/user/getmail.timer
[Unit]
Description = fetch mail from server periodically

[Timer]
OnStartupSec = 5m
OnUnitActiveSec = 5m

[Install]
WantedBy = timers.target
