# systemd/.config/systemd/user/mpd.service
[Unit]
Description = music player daemon

[Service]
ExecStart = /usr/bin/mpd --no-daemon
Restart = always

