# systemd/.config/systemd/user/mpd.timer
[Unit]
Description = start mpd

[Timer]
OnBootSec = 1h

[Install]
WantedBy = timers.target
