# systemd/.config/systemd/user/pacman_diff.timer
[Unit]
Description = report changes on pacman packages

[Timer]
OnStartupSec = 30m

[Install]
WantedBy = timers.target
