# systemd/.config/systemd/user/ted_talks.timer
[Unit]
Description = Run TED Talks service every day

[Timer]
OnCalendar = 21:05

[Install]
WantedBy = timers.target
