# systemd/.config/systemd/user/xsession.target
[Unit]
Description = xsession services
After = default.target
